# instance/services/firefox/run
#!/bin/sh

DISPLAY=:1.0 exec firefox about:blank
