# instance/services/audioproxy/run
#!/bin/sh

exec /home/user/audio-proxy.sh -l 5711
