# instance/services/api/run
#!/bin/sh

exec python3 /home/user/instance-api.py
